# infra/server/scripts/tier0-boot-inhibit.service
[Unit]
Description=Tier-0 Boot Inhibit -- set backup_inhibit until monitoring verifies
Documentation=State is UNKNOWN after reboot. Inhibit must exist until monitor-backups.sh clears it.
DefaultDependencies=no
Before=network.target
Before=cron.service

[Service]
Type=oneshot
# Create inhibit with a reboot marker so logs show WHY it was set
ExecStart=/bin/bash -c 'echo "REBOOT $(date --iso-8601=seconds) -- state UNKNOWN, awaiting Tier-0 verification" > /var/run/backup_inhibit'
# Also remove stale heartbeat -- it's from before the reboot, meaningless now
ExecStart=/bin/rm -f /var/run/tier0_heartbeat
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
